# extra/backlightd.service
[Unit]
Description=Backlight daemon
Requires=acpid.socket

[Service]
ExecStart=/usr/local/bin/backlightd

[Install]
WantedBy=multi-user.target
