# extra/backlightd.service
[Unit]
Description=Backlight daemon

[Service]
ExecStart=/usr/lib/backlightd/backlightd

[Install]
WantedBy=multi-user.target
